# srcs/requirements/django_auth/conf/crontab
# Run a script every minute
0 * * * * python3 manage.py delete_old_accounts
